# deploy/systemd-drop-ins/prometheus.service.d/tailscale-bind.conf
# Override for prometheus.service — fix boot-time race against tailscaled.
#
# Prometheus is configured (via /etc/default/prometheus → $ARGS) to bind
# on the host's tailscale IP (e.g. 100.68.41.91:9090). On reboot, systemd
# starts prometheus.service before tailscaled has brought up the
# tailscale0 interface and assigned the IP, so the bind fails with:
#
#   level=error msg="Unable to start web listener"
#   err="listen tcp 100.68.41.91:9090: bind: cannot assign requested address"
#
# The upstream unit has Restart=on-abnormal, which does NOT cover an
# exit-code-1 failure — so the service stays dead until manual restart.
# Grafana datasource then returns "connection refused" on every query
# and dashboards go blank ("No data").
#
# This override:
#   1. Adds After/Wants on tailscaled.service + network-online.target so
#      systemd at least *waits* for the tailscale daemon (interface may
#      still come up slightly later, hence #2).
#   2. Switches Restart to on-failure with a 5 s back-off, and raises
#      the start-limit burst so the service can retry for ~100 s while
#      the tailscale interface is being negotiated. Once the IP is
#      assigned the next retry binds successfully and the service stays
#      up.
#
# Install: copy into /etc/systemd/system/prometheus.service.d/, then:
#   sudo systemctl daemon-reload
#   sudo systemctl restart prometheus

[Unit]
After=tailscaled.service network-online.target
Wants=tailscaled.service network-online.target

[Service]
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=20
