# deploy/systemd-drop-ins/prometheus-node-exporter.service.d/tailscale-bind.conf
# Override for prometheus-node-exporter.service — same race as prometheus.
#
# /etc/default/prometheus-node-exporter sets:
#   ARGS="--web.listen-address=100.68.41.91:9100"
#
# On reboot, node-exporter races tailscaled the same way Prometheus
# does (see ../prometheus.service.d/tailscale-bind.conf for the full
# rationale). The upstream unit already has Restart=on-failure, but the
# default StartLimitBurst exhausts before tailscale comes up, and there
# is no After=tailscaled dependency.
#
# This override:
#   1. Adds After/Wants on tailscaled.service + network-online.target.
#   2. Raises StartLimitBurst to 30 over 120 s so the service keeps
#      retrying while tailscale brings the interface up.
#
# Install: copy into /etc/systemd/system/prometheus-node-exporter.service.d/,
# then:
#   sudo systemctl daemon-reload
#   sudo systemctl restart prometheus-node-exporter

[Unit]
After=tailscaled.service network-online.target
Wants=tailscaled.service network-online.target

[Service]
RestartSec=5
StartLimitIntervalSec=120
StartLimitBurst=30
